# prod/freebsd/rc.d/wagtail
#!/bin/sh

# PROVIDE: wagtail
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="wagtail"
rcvar="${name}_enable"

load_rc_config "$name"

: ${wagtail_enable:=NO}
: ${wagtail_user:=www}
: ${wagtail_group:=www}
: ${wagtail_command:=/usr/local/www/wagtail/prod/freebsd/scripts/start.sh}
: ${wagtail_child_pidfile:=/var/run/wagtail/child.pid}
: ${wagtail_supervisor_pidfile:=/var/run/wagtail/supervisor.pid}
: ${wagtail_sock:=/var/run/wagtail/wagtail.sock}
: ${wagtail_log:=/var/log/wagtail.log}
: ${wagtail_rundir:=/var/run/wagtail}

command="/usr/sbin/daemon"

command_args="-f -u ${wagtail_user} -p ${wagtail_child_pidfile} -P ${wagtail_supervisor_pidfile} -o ${wagtail_log} -m 3 ${wagtail_command}"

extra_commands="reload"
reload_cmd="wagtail_reload"

start_precmd="wagtail_prestart"
stop_postcmd="wagtail_poststop"

# ---- hooks ---- #

wagtail_prestart()
{
    # create runtime dir (every boot)
    install -d -o ${wagtail_user} -g ${wagtail_group} -m 755 ${wagtail_rundir}

    rm -f "${wagtail_sock}" "${wagtail_child_pidfile}" "${wagtail_supervisor_pidfile}" 2>/dev/null || true
}

wagtail_poststop()
{
    # cleanup socket after stop
    [ -S "${wagtail_sock}" ] && rm -f "${wagtail_sock}"
}

wagtail_reload()
{
    if [ -f "${wagtail_child_pidfile}" ]; then
        kill -TERM "$(cat ${wagtail_child_pidfile})" 2>/dev/null || true
    fi
}

run_rc_command "$1"
